# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

[program:php]
command=/usr/sbin/php5-fpm

[program:nginx]
command=nginx -g "daemon off;"
