# supervisord.conf
[supervisord]
nodaemon=true

[program:php]
command=/usr/sbin/php5-fpm

[program:nginx]
command=nginx -g "daemon off;"
